# package/utils/ibmtpm20tss/files/ibmtpm20tss.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 A. Pape

START=50

start() {
	export TPM_DEVICE=/dev/tpm0
	export TPM_INTERFACE_TYPE=dev
}

